# VirtualBox-bin/systemd.d/vboxautostart.service
[Unit]
Description=VirtualBox Autostart

[Service]
ExecStart=/etc/init.d/vboxautostart-service start

[Install]
WantedBy=basic.target
