# VirtualBox-bin/systemd.d/vboxweb.service
[Unit]
Description=VirtualBox Web services

[Service]
ExecStart=/etc/init.d/vboxweb-service start

[Install]
WantedBy=basic.target
