# VirtualBox-bin/systemd.d/virtualbox.service
[Unit]
Description=VirtualBox guest services

[Service]
ExecStart=/etc/init.d/vboxdrv start

[Install]
WantedBy=basic.target
